# modules/services/.config/systemd/user/time_reminder.service
[Unit]
Description=Notify test

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c 'notify-send "$(date +%%H:%%M)" " $(python3 $HOME/Main/Scripts/RAM_Script_Python/main.py random --hide-id --hide-date) (random todo)"'

[Install]
WantedBy=default.target
